# grpc/docker/grpc_server/supervisord.conf
[supervisord]
nodaemon=true

[program:grpc_server]
command=/opt/virtualenv/bin/python /opt/virtualenv/lib/python2.7/site-packages/pynms_grpc/server/pynms_grpc_server.py
