# grpc/docker/grpc_server/supervisord.conf
[supervisord]
nodaemon=true

[program:grpc_server]
command=/opt/virtualenv/bin/python /opt/pynms_server/bin/grpc_server.py
